# src/tools/startup_services/services/frc-ros_docker.service
[Unit]
Description=FRC-ROS Docker Container Launch
Wants=network-online.target
After=network-online.target

[Service]
EnvironmentFile=/etc/default/frc-ros
Type=simple
ExecStart=/bin/bash ${FRC_ROS_DIR}/src/tools/startup_services/launch/docker.bash
ExecStop=/usr/bin/docker stop ${CONT_NAME}
TimeoutStartSec=600s
TimeoutStopSec=5s

[Install]
WantedBy=multi-user.target

